# scorpion.service
[Unit]
Description=Scorpion Security Platform
After=network.target
Wants=network.target

[Service]
Type=simple
User=scorpion
Group=scorpion
WorkingDirectory=F:\Testing_Tool
Environment=NODE_ENV=production
Environment=PORT=3001
ExecStart=/usr/bin/node server/simple-web-server.js
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=scorpion

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=F:\Testing_Tool/logs F:\Testing_Tool/reports F:\Testing_Tool/results

[Install]
WantedBy=multi-user.target